# system/ruby.server.gitwiki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ruby.server.gitwiki</string>
    
    <key>WorkingDirectory</key>
    <!-- change to path to your clone of git-wiki -->
    <string>/home/smueller/Code/Wiki/git-wiki</string>
    
    <key>Program</key>
    <string>git-wiki.rb</string>
    
    <key>ProgramArguments</key>
    <array>
      <string>-e production</string>
      <string>-p 8777</string>
    </array>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>UserName</key>
    <!-- change to your username -->
    <string>zzeroo</string>
    
    <key>Debug</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
    	<key>PATH</key>
      <!--<string>/usr/local/bin:/opt/local/bin:/opt/local/sbin:/usr/local/mysql/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin:/usr/X11/bin</string>-->
      <string>/home/smueller/.rvm/gems/ruby-1.8.7-p299/bin:/home/smueller/.rvm/gems/ruby-1.8.7-p299@global/bin:/home/smueller/.rvm/rubies/ruby-1.8.7-p299/bin:/home/smueller/.rvm/bin:/home/smueller/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin</string>
    </dict>
  </dict>
</plist>
